# pi/fallback-ap/systemd/projectplant-apply-wifi.service
[Unit]
Description=ProjectPlant: Apply Wi‑Fi credentials and revert AP
After=NetworkManager.service

[Service]
Type=oneshot
ExecStart=/opt/projectplant/fallback-ap/projectplant-apply-wifi.sh
User=root
Group=root

[Install]
WantedBy=multi-user.target

